# init/icmpPing.ini
[program:icmp_clientd]
command=/usr/local/icmpPing/src/client.py -c /usr/local/icmpPing/conf/client.ini
priority=9
autostart=true
autorestart=true

[program:icmp_serverd]
command=/usr/local/icmpPing/src/server.py
priority=9
autostart=true
autorestart=true
